# pgbouncer/supervisor/pgbouncer.conf
[program:pgbouncer]
user=postgres
command=/usr/sbin/pgbouncer /etc/pgbouncer/pgbouncer.ini 
