# consul/supervisor/consul.conf
[program:consul]
environment=GOMAXPROCS="`nproc`"
command=/usr/local/bin/consul agent -join="10.8.0.1" -config-dir="/etc/consul/consul.d" -data-dir="/var/consul" 
