# consul/supervisor/consul.conf
[program:consul]
environment=GOMAXPROCS="`nproc`"
command=/usr/local/bin/consul agent -join consul.service.jaguar.consul -dc jaguar -config-dir="/etc/consul/consul.d" -data-dir="/tmp/consul" 
